# horizon/files/supervisor.conf
{%- set service = salt['pillar.get']('supervisor:server:service:'+service_name) %}
{%- set app = salt['pillar.get']('horizon:server:app:'+service.name) %}

[program:{{ service_name }}]
environment=PATH="/srv/horizon/bin"

command=/srv/horizon/sites/{{ service.name }}/gunicorn_start

stdout_logfile=/var/log/horizon/{{ service.name }}_app.log 
stderr_logfile=/var/log/horizon/{{ service.name }}_error.log

user=horizon
autostart=true
autorestart=true
